# datastax-cassandra/etc/service/cassandra/run
#!/bin/bash
set -e

set -- cassandra -f

exec "$@"
