# datastax-cassandra/etc/service/opscenter-agent/run
#!/bin/bash
# Run opscenter agent
sv -w 15 check cassandra || exit 1
exec /usr/local/opscenter/datastax-agent*/bin/datastax-agent -f
