# datastax-cassandra/etc/service/opscenter-agent/run
#!/bin/bash
# Run opscenter agent
exec /usr/local/opscenter/datastax-agent*/bin/datastax-agent -f
