# datastax-opscenter/etc/service/opscenterd/run
#!/bin/bash
exec /etc/init.d/opscenterd start
