# datastax-opscenter/etc/service/opscenterd/run
#!/bin/bash
exec /usr/share/opscenter/bin/opscenter -f
